# board/roles/smartcar_server/templates/smartcar-server.service.j2
{# Transport-aware smartcar-server unit. Forks Unit ordering + Conditions
   on smartcar_transport. #}
[Unit]
{% if smartcar_transport == "bt" %}
Description=smartcar Android Auto projection (AAW — BR/EDR + WiFi STA)
After=bluetooth.service systemd-networkd.service
Wants=bluetooth.service systemd-networkd.service
# No StartLimitBurst on bt — the bt module's pair_wait_timeout loops can
# legitimately exit after 10 min if the car never appeared, and we want
# systemd to bring us right back up so the operator can keep advertising
# indefinitely without manual `systemctl restart` after each window.
{% elif smartcar_transport == "usb" %}
Description=smartcar Android Auto projection (car mode — USB AOAP)
After=usb-mode.service
Requires=usb-mode.service
ConditionPathExists=/run/usb-car-mode
StartLimitIntervalSec=60
StartLimitBurst=5
{% else %}
Description=smartcar Android Auto projection (dev — TCP to openauto)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=5
{% endif %}

[Service]
Type=simple
Environment=RUST_LOG={{ smartcar_rust_log }}
# Optional transient runtime override written by scripts/deploy.py via
# `--runtime-args`. /run is tmpfs — the file evaporates on reboot, never
# persists. The leading '-' on EnvironmentFile makes it optional, so the
# unit starts normally when no override is set.
EnvironmentFile=-/run/smartcar-deploy/runtime.env
ExecStart={{ smartcar_binary }} --transport {{ smartcar_transport }} --bridge {{ smartcar_bridge }}{% if smartcar_extra_args %} {{ smartcar_extra_args }}{% endif %} $DEPLOY_EXTRA_ARGS
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
