# board/roles/smartcar_server/templates/smartcar-server.service.j2
[Unit]
Description=smartcar Android Auto projection (car mode — USB AOAP)
After=usb-mode.service
Requires=usb-mode.service
ConditionPathExists=/run/usb-car-mode
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
Environment=RUST_LOG={{ smartcar_rust_log }}
# Optional transient runtime override written by scripts/deploy.py via
# `--runtime-args`. /run is tmpfs — the file evaporates on reboot, never
# persists. The leading '-' on EnvironmentFile makes it optional, so the
# unit starts normally when no override is set.
EnvironmentFile=-/run/smartcar-deploy/runtime.env
ExecStart={{ smartcar_binary }} --transport usb{% if smartcar_extra_args %} {{ smartcar_extra_args }}{% endif %} $DEPLOY_EXTRA_ARGS
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
